# deploy/ssqs.service
[Unit]
Description=SSQS Next.js Website
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/ssqs/current
Environment=NODE_ENV=production
Environment=PORT=3100
Environment=HOSTNAME=127.0.0.1
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env node /var/www/ssqs/current/.next/standalone/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
